# install/inkypi.service
[Unit]
Description=InkyPi App
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=root
RuntimeDirectory=inkypi
WorkingDirectory=/run/inkypi
ExecStart=/usr/local/bin/inkypi run
Restart=on-failure
RestartSec=60
WatchdogSec=120
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal
CPUQuota=40%
MemoryHigh=250M
MemoryMax=350M
# JTN-601: Prefer inkypi as OOM victim over sshd/systemd-journald so we can
# still SSH in during memory crunch. Positive value makes the kernel more
# likely to kill inkypi first. See install/inkypi.service history.
OOMScoreAdjust=500

[Install]
WantedBy=multi-user.target